# files/supervisord.conf
[supervisord]
nodaemon=true

[program:iventoy]
directory=/iventoy
command=sh -c 'env IVENTOY_API_ALL=1 lib/iventoy'
